# todo-bot.service
[Unit]
Description=Kundalik Vazifalar Telegram Bot
After=network.target postgresql.service

[Service]
Type=simple
User=todobot
WorkingDirectory=/opt/todo_bot
EnvironmentFile=/opt/todo_bot/.env
ExecStart=/opt/todo_bot/venv/bin/python -m app.bot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
